# scraper/ahorrapy-scrape.timer
[Unit]
Description=Run AhorraPY Scraper Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
